# roles/backend/files/dev/kestrel-jleague-web.service
[Unit]
Description=BuggyAmb is a really buggy application
 
[Service]
WorkingDirectory=/var/www/html/jleague/web/build
ExecStart=/usr/bin/dotnet/dotnet /var/www/html/jleague/web/build/JLeague.Web.dll
Restart=always
# Restart service after 10 seconds if the dotnet service crashes:
RestartSec=10
KillSignal=SIGINT
SyslogIdentifier=buggyamb-identifier
User=root
Environment=ASPNETCORE_ENVIRONMENT=Development
Environment=DOTNET_PRINT_TELEMETRY_MESSAGE=false
Environment=ASPNETCORE_URLS=http://localhost:4000
EnvironmentFile=/var/www/html/jleague/web/backend.env
#Environment=JL_WEB_ConnectionStrings__DefaultConnection=Server=jleague-rds.cd3mtzp0u8vb.ap-northeast-1.rds.amazonaws.com;User Id=jleague-Dev;Database=jleague_prod;Password=R6gukbQ0Qa7PBFkF;Max Pool Size=500;
#Environment=JL_WEB_IronPdf__LicenseKey=IRONPDF-1141360AC3-167450-CE4561-A5F334F237-985E4C69-UExE509D7DE0DC59D8-IRO200810.9494.74130.PRO.1DEV.1YR.SUPPORTED.UNTIL.11.AUG.2021
#Environment=JL_WEB_Uploads__BasePath=/var/www/html/jleague/web/Content
#Environment=JL_WEB_Uploads__ClubLogosSubPath=/clubLogos/
#Environment=JL_WEB_Uploads__UploadFilesPath=/var/www/html/jleague/web/uploads
#Environment=JL_WEB_Hangfire__Enabled=true
#Environment=JL_WEB_Deepl__Url=https://api.deepl.com
#Environment=JL_WEB_Deepl__AuthKey=29da9471-e3a2-d89c-b4a7-06a24d6b0ad3
#Environment=JL_WEB_PdfFocus__Serial=70056334650
#Environment=JL_WEB_PdfMetamorphosis__Serial=80090879749
#Environment=JL_WEB_CloudConvert__ApiKey=eyJ0eXAiOiJKV1QiLCJhbGciOiJSUzI1NiJ9.eyJhdWQiOiIxIiwianRpIjoiZmYxMjg1MzA0MzZkYWMyNjA4ZjZlOWQxNmJkNWQ5NWQxNzFiODZhMzdhMjM1NDc5YTU2ZDYwY2I0YjdjNDE4OTliZmI3Njc5NGM1MTJiNzUiLCJpYXQiOjE2NDkxOTI5NTAuMjU5Mzc1LCJuYmYiOjE2NDkxOTI5NTAuMjU5Mzc3LCJleHAiOjQ4MDQ4NjY1NTAuMjQ3MDksInN1YiI6IjU3MjA3ODU3Iiwic2NvcGVzIjpbInVzZXIucmVhZCIsInVzZXIud3JpdGUiLCJ0YXNrLnJlYWQiLCJ0YXNrLndyaXRlIiwid2ViaG9vay5yZWFkIiwid2ViaG9vay53cml0ZSIsInByZXNldC5yZWFkIiwicHJlc2V0LndyaXRlIl19.mP2MUNziby4soQdTNSryocancrr0OUo4VYV_MP2yFAzqa5BHMhl-BYnnSlVDmOakFsGDmFuZZQ90FVt6c1BrUhJOGkpBTOYwixcVB_FUNqA_2lOtbafeHUO6H-YiWBt3fq3oixjMHPZszbfdYDUWKtrTrw02vWJGLoCP8b4OI2sDzqhVnfaATXHrGKnxVJvINysFCY3l5--f4WJoQayq7Sv5ZCn3d9rdYxEErw3nCg1zhmkmAitoqpB_mQElL_YFBW5-II6GnU911E-tHNWL90JTTGWmDzIvpWUyhFTnNmViRfrPH89sy6p8p8tkhtk89jUtk5A-q1_tubGmd1Qm-h5NY58UQWTEKgrSg6C8spzsDdeIh-OxE9x90SnVf4PFzUZK8uau4M_QMLE_bUfIGsJ3i8MiAN9WKW6EDiP8lY29ek1Ss8VfOtfc-n30VJM2VavjK1JFcS3_SsMPuBY5xporY67cuLYUeB1Ds47JaFudyG4SJRTo5EaCQSFoqDVK56Pn_tsSznkQjBnL7LclDtxOpnDq-lC0yKgePkfpxL6VzcdBYpXCB6z8fRysBdSvpuGv4o5MZuitsBLTpIcM3CvivE-YYfyutGg4TqugsSW6ZrQDs53a9CWgm4TyqoKG5RfI3VuJsjrzw-Hm8GytNMWgtGWjKN7W1-dC2dmnLTo
#Environment=JL_WEB_CloudConvert__IsSandBox=false
#Environment=JL_WEB_EnabledFeatures__AutoMigrations=true
 
[Install]
WantedBy=multi-user.target