# roles/backend/files/dev/kestrel-jleague-admin.service
[Unit]
Description=Example .NET Web API App running on Ubuntu

[Service]
WorkingDirectory=/var/www/html/jleague/web/admin
ExecStart=/usr/bin/dotnet/dotnet /var/www/html/jleague/web/admin/JLeague.Manager.dll
Restart=always
# Restart service after 10 seconds if the dotnet service crashes:
RestartSec=30
KillSignal=SIGINT
SyslogIdentifier=dotnet-example
User=root
Environment=ASPNETCORE_ENVIRONMENT=Publish
Environment=DOTNET_PRINT_TELEMETRY_MESSAGE=false
EnvironmentFile=/var/www/html/jleague/web/admin.env
#Environment=JL_MGR_Uploads__BasePath=/var/www/html/jleague/web/Content
#Environment=JL_MGR_Uploads__ClubLogosSubPath=/clubLogos/
#Environment=JL_MGR_Uploads__UploadFilesPath=/var/www/html/jleague/web/uploads
#Environment=JL_MGR_FilePrefixes__AdminPrefix=admin/
#Environment=JL_MGR_Email__SmtpHost=
#Environment=JL_MGR_Email__SmtpPort=
#Environment=JL_MGR_Email__UseSsl=
#Environment=JL_MGR_Email__UserName=JLeague
#Environment=JL_MGR_Email__Email=
#Environment=JL_MGR_Email__Password=
#Environment=JL_MGR_Email__SupportEmail=jleaguesupport@corvus.london
#Environment=JL_MGR_Deepl__Url=https://api.deepl.com
#Environment=JL_MGR_Deepl__AuthKey=29da9471-e3a2-d89c-b4a7-06a24d6b0ad3
#Environment=JL_MGR_PdfFocus__Serial=70056334650
#Environment=JL_MGR_CloudConvert__ApiKey=eyJ0eXAiOiJKV1QiLCJhbGciOiJSUzI1NiJ9.eyJhdWQiOiIxIiwianRpIjoiZmYxMjg1MzA0MzZkYWMyNjA4ZjZlOWQxNmJkNWQ5NWQxNzFiODZhMzdhMjM1NDc5YTU2ZDYwY2I0YjdjNDE4OTliZmI3Njc5NGM1MTJiNzUiLCJpYXQiOjE2NDkxOTI5NTAuMjU5Mzc1LCJuYmYiOjE2NDkxOTI5NTAuMjU5Mzc3LCJleHAiOjQ4MDQ4NjY1NTAuMjQ3MDksInN1YiI6IjU3MjA3ODU3Iiwic2NvcGVzIjpbInVzZXIucmVhZCIsInVzZXIud3JpdGUiLCJ0YXNrLnJlYWQiLCJ0YXNrLndyaXRlIiwid2ViaG9vay5yZWFkIiwid2ViaG9vay53cml0ZSIsInByZXNldC5yZWFkIiwicHJlc2V0LndyaXRlIl19.mP2MUNziby4soQdTNSryocancrr0OUo4VYV_MP2yFAzqa5BHMhl-BYnnSlVDmOakFsGDmFuZZQ90FVt6c1BrUhJOGkpBTOYwixcVB_FUNqA_2lOtbafeHUO6H-YiWBt3fq3oixjMHPZszbfdYDUWKtrTrw02vWJGLoCP8b4OI2sDzqhVnfaATXHrGKnxVJvINysFCY3l5--f4WJoQayq7Sv5ZCn3d9rdYxEErw3nCg1zhmkmAitoqpB_mQElL_YFBW5-II6GnU911E-tHNWL90JTTGWmDzIvpWUyhFTnNmViRfrPH89sy6p8p8tkhtk89jUtk5A-q1_tubGmd1Qm-h5NY58UQWTEKgrSg6C8spzsDdeIh-OxE9x90SnVf4PFzUZK8uau4M_QMLE_bUfIGsJ3i8MiAN9WKW6EDiP8lY29ek1Ss8VfOtfc-n30VJM2VavjK1JFcS3_SsMPuBY5xporY67cuLYUeB1Ds47JaFudyG4SJRTo5EaCQSFoqDVK56Pn_tsSznkQjBnL7LclDtxOpnDq-lC0yKgePkfpxL6VzcdBYpXCB6z8fRysBdSvpuGv4o5MZuitsBLTpIcM3CvivE-YYfyutGg4TqugsSW6ZrQDs53a9CWgm4TyqoKG5RfI3VuJsjrzw-Hm8GytNMWgtGWjKN7W1-dC2dmnLTo
#Environment=JL_MGR_CloudConvert__IsSandBox=false

[Install]
WantedBy=multi-user.target